# ansible/roles/base/files/nextcloud-sync.service
[Unit]
Description=nextcloudcmd sync job

[Service]
Type=oneshot
ExecStart=/home/cameron/.local/bin/nextcloud-sync
