# ansible/roles/base/files/vdirsyncer.service
[Unit]
Description=vdirsyncer sync job

[Service]
Type=oneshot
ExecStart=/usr/bin/vdirsyncer sync
